# systemd/gallery-portal.service
[Unit]
Description=Zen Gallery Portal
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/wizzard/gpu-ai-hub/src
ExecStart=/home/wizzard/ai/gallery/venv/bin/uvicorn gallery_portal:app --host 0.0.0.0 --port 8009
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
