# systemd/tts-portal.service
[Unit]
Description=Zen Voice Pro (Piper & XTTS v2)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/wizzard/ai/tts/xtts
ExecStart=/home/wizzard/ai/tts/xtts/xtts-env/bin/python -m uvicorn xtts_portal:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=COQUI_TOS_AGREED=1

[Install]
WantedBy=default.target
