# systemd/zen-hub.service
[Unit]
Description=Zen AI Hub Landing Page
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/wizzard/ai/comfyui/workflow-launcher
ExecStart=/home/wizzard/ai/comfyui/venv/bin/python -m http.server 8191 --bind 0.0.0.0
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
